# system/code_learning_club.crontab
## The schedule
* * * * * php /papac-and-co/app-production/artisan schedule:run
